# updatetool/rc.local
echo 2048 > /proc/sys/vm/min_free_kbytes
./etc/rc.d/mac.sh
./etc/rc.d/ip.sh
#temp_ip

if [ -f "/nand/UpFiles/update.sh" ]; then
  chmod +x /nand/UpFiles/update.sh
  cd /nand/UpFiles
  ./update.sh
  rm /nand/UpFiles/*
fi

syslogd -l 7 -O /nand/log
mkfifo /dev/shm/null        
tail -f /dev/shm/null>> /dev/null &
cjmain all&
